# systemd/handoff-emea.service
[Unit]
Description=Handoff snapshot EMEA

[Service]
ExecStart=/usr/bin/node /opt/handoff/handoff_snapshot.mjs emea
EnvironmentFile=/opt/handoff/.env
