# systemd/handoff-us.service
[Unit]
Description=Handoff snapshot US

[Service]
ExecStart=/usr/bin/node /opt/handoff/handoff_snapshot.mjs us
EnvironmentFile=/opt/handoff/.env
